# mini-control.service
[Unit]
Description=Mac Mini Control Panel
After=network.target

[Service]
User=ludovic
WorkingDirectory=/home/ludovic/mini-control
ExecStart=/home/ludovic/.local/bin/gunicorn --workers 1 --threads 2 --preload -b 0.0.0.0:5000 app:app
Restart=on-failure
Environment=PATH=/home/ludovic/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
MemoryMax=150M

[Install]
WantedBy=multi-user.target
